# extra/asahi-battery-threshold.service
[Unit]
Description=asahi-battery-threshold

[Service]
ExecStart=/usr/bin/asahi-battery-threshold

[Install]
WantedBy=multi-user.target
